# docker/crontab
# Phase 1: nightly DB backup at 03:00. Later phases add pipeline/notify jobs here.
0 3 * * * skannonser db backup
