# docker/crontab
# Nightly DB backup at 03:00 UTC, keep newest 30.
0 3 * * * skannonser db backup --keep 30
